# scripts/debian/postinst
#!/bin/sh
if [ "$(pidof systemd)" ]; then
   systemctl preset DCore &>/dev/null
fi
