# scripts/debian/prerm
#!/bin/sh
if [ "$(pidof systemd)" ]; then
   systemctl --global disable DCore &>/dev/null
fi
