# scripts/debian/prerm
#!/bin/sh
if [ "$(pidof systemd)" ]; then
   systemctl --no-reload disable --now DCore &>/dev/null
fi
